# ansible/roles/cloudflare-ddns/templates/cloudflare-ddns.service.j2
[Unit]
Description=Cloudflare DDNS update
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/cloudflare-ddns/env
ExecStart=/usr/local/bin/cloudflare-ddns
